# systemd/galilB.service
[Unit]
Description=M2FS Galil B

[Install]
WantedBy=director.service
#Also=

[Service]
Type=simple
ExecStart=/M2FS-Control/galilAgent.py -B
RestartSec=100ms
Restart=always
KillMode=control-group
KillSignal=SIGTERM
SendSIGKILL=yes